# mlapi.conf
[program:nginx]
command=/usr/sbin/nginx

[program:gunicorn]
command=/miniconda/bin/gunicorn app:app -b localhost:8000
directory=/var/www/html
